# files/Swift-Arch-ww/systemd/paccacheu.service
[Unit]
Description=Remove uninstalled cached package files

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -ruk0
